# scripts/code-watchdog.service
[Unit]
Description=code-watchdog — auto-resume dropped Claude sessions
After=default.target

[Service]
Type=oneshot
ExecStart=/home/brose/bin/code-watchdog.sh
StandardOutput=journal
