# ops/claude-jsonl-heal.timer
[Unit]
Description=Run claude-jsonl-heal hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=5min
Persistent=true

[Install]
WantedBy=default.target
